# mrm_fedora_update.service
[Unit]
Description=MRM Fedora Startup

[Service]
ExecStart=/usr/bin/fedora_update.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
